# backend/system/clisonix_daily.service
[Unit]
Description=Clisonix Daily Report Sender
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/backend/system
ExecStart=/usr/bin/python3 daily_scheduler.py
Restart=always
RestartSec=10
Environment="SMTP_HOST=smtp.gmail.com"
Environment="SMTP_PORT=587"
Environment="SMTP_USER=amati.bau@gmail.com"
Environment="SMTP_PASS=APPPASSWORD16"
Environment="REPORT_TO=amati.bau@gmail.com"
Environment="REPORT_FROM=Clisonix Reports <amati.bau@gmail.com>"

[Install]
WantedBy=multi-user.target
